# config/messenger-worker.ini
[program:messenger-consume]
command=php /appli/questionnaire/current/bin/console messenger:consume async --time-limit=3600
user=nginx
numprocs=2
autostart=true
autorestart=true
process_name=%(program_name)s_%(process_num)02d